# supervisord/etc/supervisor/conf.d/firstboot.conf
[program:firstboot]
priority=1
autorestart=false
command=/bin/sh /etc/supervisor/firstboot
