# ts/etc/supervisor/conf.d/fail2ban.conf
[program:fail2ban]
command=/usr/bin/fail2ban-server -f -s /var/run/fail2ban/fail2ban.sock
