# ts/etc/supervisor/conf.d/xrdp-sesman.conf
[program:xrdp-sesman]
command=bash -c "rm -f /var/run/xrdp/*.pid && /usr/local/sbin/xrdp-sesman --nodaemon"
user=root
autorestart=true
priority=1010
